# playbooks/roles/jenkins_common/templates/etc/systemd/system/jenkins.service.j2
[Unit]
Description=Jenkins

[Service]
EnvironmentFile=/etc/default/jenkins
User=jenkins
Group=jenkins
ExecStart=/usr/bin/java ${JAVA_ARGS} -jar ${JENKINS_WAR} ${JENKINS_ARGS}

[Install]
WantedBy=multi-user.target
